# crontab
*	13	*	*	*	python /app/redmine-tidy.py

